# supervisord.conf
[supervisord]

[program:cat]
command=/bin/cat
